# luci-0.10/applications/luci-astral/root/etc/uci-defaults/luci-astral
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@astral[-1]
	add ucitrack astral
	set ucitrack.@astral[-1].init=astral
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
